# net/busybeed/pkg/busybeed.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/busybeed"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
